# trant/.config/systemd/user/trantbox-games.service
[Unit]
Description=TranThangBin's games container
PartOf=trantbox.target

[Service]
Type=oneshot
ExecStart=/usr/bin/podman container start trantbox-games
ExecStop=/usr/bin/podman container stop trantbox-games
RemainAfterExit=yes
