# trant/.config/systemd/user/trantbox@.target
[Unit]
Description=TranThangBin's %i container
Requires=podman-container@trantbox-%i.service
BindsTo=podman-container@trantbox-%i.service
