# stow/startup/dot-config/systemd/user/trantbox-services-easyeffects.service
# Reference from https://github.com/nix-community/home-manager/blob/release-25.05/modules/services/easyeffects.nix
[Unit]
Description=Easyeffects daemon
PartOf=graphical-session.target
After=graphical-session.target

Requires=distrobox-container@trantbox-services.service
PartOf=distrobox-container@trantbox-services.service
Conflicts=easyeffects.service

[Service]
ExecCondition=distrobox-enter trantbox-services -- sh -c 'command -v easyeffects &>/dev/null'
ExecStart=distrobox-enter trantbox-services -- easyeffects --hide-window --service-mode
ExecStop=distrobox-enter trantbox-services -- easyeffects --quit
KillMode=mixed
Restart=on-failure
RestartSec=5
TimeoutStopSec=10
