# etc/systemd/system/systemd-hibernate.service.d/50_custom.conf
[Unit]
Wants=nvidia-hibernate.service nvidia-resume.service
