# etc/systemd/system/systemd-suspend.service.d/50_custom.conf
[Unit]
Wants=nvidia-resume.service nvidia-suspend.service
